# alarm_player.service
[Unit]
Description=BellNews Alarm Player Service
After=network.target sound.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/bellapp
ExecStart=/usr/bin/python3 /root/bellapp/alarm_player.py
Restart=always
RestartSec=10
StartLimitInterval=0
StandardOutput=journal
StandardError=journal

# Watchdog to detect hung processes
WatchdogSec=300

# Resource limits to prevent runaway processes
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
